# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:bot]
command=/app/bot
