# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:bot]
command=/app/bot

[program:api]
command=/app/api
